# supervisord.conf
[supervisord]
nodaemon=true

[program:logstash-forwarder]
command=/opt/lumberjack/bin/lumberjack -config /etc/logstash-forwarder/logstash-forwarder.conf

[program:nginx]
command=nginx
